# data/supervisor.conf
[program:texttospeech]
command = gunicorn -b 0.0.0.0:12432 app:app
directory = /opt/loxberry/data/plugins/texttospeech/scripts
user = root
environment=HOME="/root", USER="root"
